# php/56/conf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:fpm]
command=/usr/sbin/php-fpm5 -R --nodaemonize

[program:sshd]
command=/usr/sbin/sshd -D